# scripts/com.cdcoonce.portfolio-metrics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly refresh of the portfolio's featured-card KPIs.
  Runs scripts/collect-metrics.py, which recomputes metrics.json from the
  source repos and commits/pushes it. Pages rebuilds on push to master.

  Activation is documented in scripts/metrics-cron.md. Do not enable until the
  cockpit work is merged to master (that's the branch production builds from).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cdcoonce.portfolio-metrics</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/cdcoonce/Developer/GitHub/PortfolioWebsite/scripts/collect-metrics.py</string>
    <string>--write</string>
    <string>--commit</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/cdcoonce/Developer/GitHub/PortfolioWebsite</string>

  <!-- Mondays at 08:30 local. Numbers change slowly; weekly is plenty. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/cdcoonce/Library/Logs/portfolio-metrics.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/cdcoonce/Library/Logs/portfolio-metrics.log</string>
</dict>
</plist>
